# docker/sails/etc/supervisor/conf.d/sails.conf
[program:sails]
directory=/server
command=/usr/local/bin/sails lift
stdout_logfile=/tmp/sails.log
